# yocto-build-env/s6-overlay/s6-rc.d/lock-manager/finish
#!/usr/bin/env bash

rm -f /tmp/balena/updates.lock
